# _examp_see/alt_bash.d/etc_bash.d/profile.d/cvs.sh
# default system cvsroot
#export CVSROOT=/home/cvs/cvsroot
# default transport for CVS RSH protocol
export CVS_RSH=ssh
