# _examp_see/alt_bash.d/etc_bash.d/profile.d/mc.sh
# define aliases for zsh
[ -n "${ZSH_VERSION}" ] || return 0
alias mc='. /usr/lib/mc/mc-wrapper.sh -r'
