# _examp_see/alt_bash.d/etc_bash.d/profile.d/systemd.sh
export SYSTEMD_PAGER="/usr/bin/less -FR"
